# workspace/stock-tracker/c2506-monitor.service
[Unit]
Description=C2506 Corn Futures Monitor - Hourly Report
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /home/admin/.openclaw/workspace/stock-tracker/c2506_monitor.py --mode hourly --config /home/admin/.openclaw/workspace/stock-tracker/c2506_config.json
WorkingDirectory=/home/admin/.openclaw/workspace/stock-tracker
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
