# workspace/stock-tracker/c2506-monitor-daily.timer
[Unit]
Description=C2506 Corn Futures Monitor - Daily Summary Timer
Documentation=man:systemd.timer(5)

[Timer]
# Run at 15:30 on weekdays (after market close at 15:00)
OnCalendar=*-*1,2,3,4,5 15:30:00
Persistent=true
AccuracySec=1min

[Install]
WantedBy=timers.target
